# templates/default/upstart.nsqadmin.conf.erb
description "nsqadmin"

start on runlevel [2345]
stop on runlevel [!2345]

env PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/sbin:/bin

respawn

script
  /usr/local/bin/nsqadmin --lookupd-http-address=<%= node["nsq"]["nsqd"]["lookupd_http_address"].first %> --http-address=<%= node["nsq"]["nsqadmin"]["http_address"] %> 2>&1
end script
